# systemd/mac-hid-setup.service
[Unit]
Description=mac.hid — USB HID Gadget Setup
Documentation=https://github.com/gargvasu/zerobridge
After=sysinit.target
Before=pi-agent.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/multi-hid-setup.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mac-hid-setup

[Install]
WantedBy=multi-user.target